# cli/sysadmin/src/main/scripts/storage-quarantine-cleanup.cron
# remove quarantined storage disks
#

*/5 * * * * root (date --iso-8601=seconds --utc; /usr/bin/stratus-storage-quarantine --clean) >>/var/log/stratuslab/storage-quarantine-cleanup.log 2>&1
